# configure/playbooks/roles/docker/templates/docker.service.j2
.include /lib/systemd/system/{{ docker.service_name }}.service

[Service]
ExecStart=
ExecStart=/usr/bin/dockerd
